# deploy/douyin.service.example
[Unit]
Description=Douyin subscription downloader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=douyin
Group=douyin
WorkingDirectory=/opt/douyin
EnvironmentFile=-/opt/douyin/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/douyin/.venv/bin/python /opt/douyin/main.py --host 127.0.0.1 --port 8765
Restart=always
RestartSec=0
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
